# deploy/robocompanion.service
[Unit]
Description=RoboCompanion v2 Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/robocompanion
ExecStart=/opt/robocompanion/venv/bin/python bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
